# bin/ledscape.service
#
# Copy this file to /etc/systemd/system/ and run:
#   systemctl enable ledscape.service
# to enable LEDscape to start on each boot.
#

[Unit]
Description=BeagleBone LED matrix driver

[Service]
WorkingDirectory=/opt/LEDscape/
ExecStart=/opt/LEDscape/bin/run-ledscape
KillMode=process

[Install]
WantedBy=multi-user.target
